# root/etc/services.d/icecast/run
#!/usr/bin/with-contenv bash

exec \
    /usr/bin/icecast -c /etc/icecast.xml
